# mica/files/etc/uci-defaults/50_mica_avrdude
#!/bin/sh

[ "$(uci get system.@system[-1].mica_firstboot)" = "1" ] && return 0

. /lib/ramips.sh
board=$(ramips_board_name)

[ "$board" = "mica7688" ] || {
	return 0
}

